# packages/nxe-devcontainer/src/lib/config/supervisord.conf
[inet_http_server]
port=127.0.0.1:30055

[supervisorctl]
serverurl=http://127.0.0.1:30055

[program:nxe-k3s-registry-proxy]
command = caddy reverse-proxy --from 127.0.0.1:5443 --to nxe-k3s-registry:5000 --internal-certs --disable-redirects
depends_on = caddy-service

[program:caddy-service]
command = caddy run
